# deploy/imac/com.coldwaterkim.pocketbase.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.coldwaterkim.pocketbase</string>

    <key>UserName</key>
    <string>kimchansu</string>

    <key>GroupName</key>
    <string>staff</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/kimchansu</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>LC_ALL</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/kimchansu/.local/share/coldwaterkim/home-server/bin/pocketbase</string>
        <string>serve</string>
        <string>--http=127.0.0.1:8090</string>
        <string>--dir=/Users/kimchansu/.local/share/coldwaterkim/home-server/pb_data</string>
        <string>--migrationsDir=/Users/kimchansu/.local/share/coldwaterkim/home-server/pb_migrations</string>
        <string>--origins=https://coldwaterkim.com,https://www.coldwaterkim.com,http://127.0.0.1:5173,http://localhost:5173</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kimchansu/.local/share/coldwaterkim/home-server</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/kimchansu/Library/Logs/coldwaterkim-pocketbase.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kimchansu/Library/Logs/coldwaterkim-pocketbase.err.log</string>
</dict>
</plist>
