# deploy/imac/com.coldwaterkim.caddy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.coldwaterkim.caddy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/caddy</string>
        <string>run</string>
        <string>--config</string>
        <string>/Users/kimchansu/.local/share/coldwaterkim/home-server/Caddyfile</string>
        <string>--adapter</string>
        <string>caddyfile</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kimchansu/.local/share/coldwaterkim/home-server</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/kimchansu</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/kimchansu/Library/Logs/coldwaterkim-caddy.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kimchansu/Library/Logs/coldwaterkim-caddy.err.log</string>
</dict>
</plist>
